# zfs/local.d/zfs-relatime.start
#!/bin/bash

# Disable relatime=temporary on pool rootfs
ROOTFS=$(zpool get -o value -H bootfs)

zfs set relatime=off $ROOTFS
